# amazon-efs-mount-watchdog.initd
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/bin/$SVCNAME"

depend() {
	need net
	use dns
}
